# pkg/validator/quadlet/testdata/err.container
## assert-error RequiredKey Container Image 0 0
## assert-error RequiredKey Container Rootfs 0 0
[Container]
